# scripts/petro-kiosk.service
# Petro exhibit kiosk. Installed by pi-setup.sh (which fills in __USER__ /
# __REPO__). Deliberately NOT enabled by default — start it by hand while
# testing, `systemctl enable` it once you want autostart at boot.
[Unit]
Description=Petro kiosk (cage + chromium)
After=network-online.target seatd.service
Wants=network-online.target seatd.service

[Service]
Type=simple
User=__USER__
# seatd (not logind) provides the seat for this headless service, so cage can
# grab the GPU + input without a login session sitting on a VT. The user must
# be in the group that owns /run/seatd.sock (video, on Pi OS).
SupplementaryGroups=video render input
Environment=LIBSEAT_BACKEND=seatd
# cage/Wayland needs a writable runtime dir for its socket; systemd creates
# /run/petro-kiosk owned by the service user and cleans it up on stop.
RuntimeDirectory=petro-kiosk
Environment=XDG_RUNTIME_DIR=/run/petro-kiosk
ExecStart=__REPO__/scripts/pi-kiosk.sh
# After cage is up, software-replug the WaveShare panel to hide the startup
# cursor. `+` runs it as root (needed for sysfs) despite User= above; `-` makes
# any failure non-fatal so a missing/odd panel can never take down the kiosk.
ExecStartPost=+-__REPO__/scripts/hide-cursor-replug.sh
# If we booted onto the printer's WiFi Direct AP, (re)create the CUPS queue the
# app prints to — its IP can change per session. `+` runs as root (lpadmin); `-`
# keeps a missing/absent printer from taking down the kiosk. No-op off that AP.
ExecStartPost=+-__REPO__/scripts/install-printer-cups.sh
# On exit you land on a TTY to maintain things; don't respawn the browser.
Restart=no

[Install]
WantedBy=multi-user.target
